# github-watchdog.service
[Unit]
Description=GitHub WatchDog

[Service]
Type=simple
ExecStart=/usr/local/bin/github-watchdog.sh

[Install]
WantedBy=default.target
